# templates/visa-base/system/etc/supervisor/conf.d/dbusd.conf
[program:dbus]
command 			= /usr/bin/dbus-daemon --nofork --nopidfile --system
priority 			= 400
autorestart 	= true
