# templates/visa-base/system/etc/supervisor/conf.d/node-exporter.conf
[program:node-exporter]
command 			= /usr/local/bin/node_exporter --web.listen-address=":9100"
priority 			= 400
autorestart 	    = true