# templates/visa-base/system/etc/supervisor/conf.d/ntpd.conf
[program:ntpd]
command=/usr/sbin/ntpd -u ntp:ntp -p /var/run/ntpd.pid -g -n
autorestart=true
user=root
priority=400