# templates/visa-base/system/etc/supervisor/conf.d/process-exporter.conf
[program:process-exporter]
command 			= /usr/local/bin/process-exporter -config.path /etc/process-exporter.yml -web.listen-address ":9200"
priority 			= 400
autorestart 	    = true