# templates/visa-apps/applications/jupyter/jupyter-lab/system/etc/supervisor/conf.d/jupyterlab.conf
[program:jupyterlab]
command=/opt/visa/jupyter/jupyterlab.sh
user=root
priority=400
autorestart=true
autostart=true
