# templates/visa-base/system/etc/xdg/autostart/welcome.desktop
[Desktop Entry]
Type=Application
Name=Welcome VISA
Exec=welcome-visa
NoDisplay=false
Hidden=false